# cloud_setup/custom_firewall.service
[Unit]
After=network.service

[Service]
ExecStart=/etc/firewallsetup/firewall

[Install]
WantedBy=default.target
